# microservices/conf/supervisor.conf
[program:app-ossp]
command = /usr/local/bin/ossp --ini /deploy/ossp.ini --py-auto-reload 1
